# 8.1/config/supervisor/supervisord-php-fpm.conf
[program:php-fpm]
command = /usr/local/sbin/php-fpm
stdout_logfile = /var/log/supervisor/php-fpm-stdout
stderr_logfile = /var/log/supervisor/php-fpm-stderr
