# 8.4/config/supervisor/supervisord-sshd.conf
[program:sshd]
command = /usr/sbin/sshd -D
stdout_logfile = /var/log/supervisor/sshd-stdout
stderr_logfile = /var/log/supervisor/sshd-stderr
