# 8.4/config/supervisor/supervisord-crond.conf
[program:cron]
# Cron will only log to syslog and nothing else...
command = /usr/sbin/cron -f
stdout_logfile = /var/log/supervisor/cron-stdout
stderr_logfile = /var/log/supervisor/cron-stderr
