# tests/.docksal/services/cli/crontab
* * * * * echo "The current date is $(date)" >> /tmp/date.txt
